# deploy/systemd/polymarket-cos-sync.service
[Unit]
Description=Polymarket COS sync daemon
After=network-online.target polymarket-materialize.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/polymarket-agent/current
Environment=PYTHONUNBUFFERED=1
Environment=POLYMARKET_ENV_FILE=/etc/polymarket-agent.env
ExecStart=/opt/polymarket-agent/current/deploy/run_cos_sync.sh
Restart=always
RestartSec=15
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
